# build/etc/init.d/firstboot
#!/bin/bash
if [ true != "$INIT_D_SCRIPT_SOURCED" ] ; then
    set "$0" "$@"; INIT_D_SCRIPT_SOURCED=true . /lib/init/init-d-script
fi
### BEGIN INIT INFO
# Provides:          skeleton
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Example initscript
# Description:       Run firstboot tasks (install puppet, enable it, attempt to converge)
### END INIT INFO

source /etc/lsb-release
PUPPET_DEB="puppetlabs-release-pc1-utopic.deb"
PUPPET_REPO_URL="https://apt.puppetlabs.com/$PUPPET_DEB"

echo "Installing wget"
apt-get install -yq wget

export DEBIAN_FRONTEND=noninteractive
echo "Downloading and installing Puppet for codename ${DISTRIB_CODENAME}"
wget -nv "$PUPPET_REPO_URL"								&& \
dpkg -i "$PUPPET_DEB"	&& \
rm "$PUPPET_DEB"			&& \
apt-get update -yq										&& \
apt-get install -yq puppet-agent

if [ $? -ne 0 ]; then
	echo "Puppet failed to install?"
	exit 1
fi

echo "Puppet installed. Starting puppet agent and contacting puppet server to attempt to bootstrap."
sed -i '' 's/START=no/START=yes/' /etc/default/puppet
puppet agent --enable
update-rc.d puppet defaults
invoke-rc.d puppet start

puppet agent -t
times=0
while [ $? -ne 0 ]; do
	times=$(($times+1))
	if [ $times -eq 5 ]; then
		echo "Node seemingly not converging?"
		break
	fi
	puppet agent -t
done

echo "Done, disabling firstboot script"
update-rc.d firstboot disable
