# build/etc/init.d/firstboot
#!/bin/bash
if [ true != "$INIT_D_SCRIPT_SOURCED" ] ; then
    set "$0" "$@"; INIT_D_SCRIPT_SOURCED=true . /lib/init/init-d-script
fi
### BEGIN INIT INFO
# Provides:          skeleton
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Firstboot tasks
# Description:       Run firstboot script (denoted by the FIRSTBOOT environment variable)
### END INIT INFO

[ -f /etc/default/firstboot ] && . /etc/default/firstboot

if [ $FIRSTBOOT != "" ]; then
	logger -t firstboot -p syslog.info "Found firstboot script $FIRSTBOOT, downloading and executing"
	curl $FIRSTBOOT | sh
else
	logger -t firstboot -p syslog.info "No firstboot script found, moving on..."
fi

logger -t firstboot -p syslog.info "Disabling firstboot"
update-rc.d firstboot disable
